# llmproxy.service
[Unit]
Description=LLM Proxy - TEI-compatible rerank and OpenAI-compatible LLM/embeddings proxy
After=network.target llama-server.service llama-reranker.service llama-embeddings.service

[Service]
Type=simple
User=noname
Group=noname
WorkingDirectory=/src/llmproxy

# LLM server (chat/completions/models)
Environment="LLMPROXY_OAILLM_BASE_URL=http://127.0.0.1:8080"
Environment="LLMPROXY_OAILLM_API_KEY="

# Reranker server (TEI rerank)
Environment="LLMPROXY_TEIRERANKER_BASE_URL=http://127.0.0.1:8082"
Environment="LLMPROXY_TEIRERANKER_API_KEY="

# Embeddings server (dedicated embeddings instance)
Environment="LLMPROXY_OAIEMBEDDINGS_BASE_URL=http://127.0.0.1:8081"
Environment="LLMPROXY_OAIEMBEDDINGS_API_KEY="

# Proxy configuration
Environment=LLMPROXY_HOST=0.0.0.0
Environment=LLMPROXY_PORT=4001
Environment=LLMPROXY_API_KEY=
Environment=LLMPROXY_LOG_LEVEL=debug

# Lock configuration (optional)
Environment=LLMPROXY_LOCK_CONFIG=/etc/llmproxy/config.yaml

# Lock script (optional) - supports Python (.py) or shell (.sh/.bash)
Environment=LLMPROXY_LOCK_SCRIPT=/etc/llmproxy/hook.py

# Backend timeouts (in seconds)
Environment=LLMPROXY_OAILLM_TIMEOUT=30
Environment=LLMPROXY_OAILLM_READ_TIMEOUT=90
Environment=LLMPROXY_TEIRERANKER_TIMEOUT=60
Environment=LLMPROXY_TEIRERANKER_READ_TIMEOUT=120
Environment=LLMPROXY_OAIEMBEDDINGS_TIMEOUT=30
Environment=LLMPROXY_OAIEMBEDDINGS_READ_TIMEOUT=60

# Setup and start
ExecStartPre=/home/noname/.local/bin/uv sync
ExecStart=/src/llmproxy/.venv/bin/python -m src.llmproxy.main
Restart=on-failure
TimeoutStopSec=90s
RestartSec=5

[Install]
WantedBy=multi-user.target
