# sample/llmproxy.service
[Unit]
Description=LLM Proxy - TEI-compatible rerank and OpenAI-compatible LLM/embeddings proxy
After=network.target llama-server.service llama-reranker.service llama-embeddings.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/src/llmproxy

# Setup and start
Environment=UV_PROJECT_ENVIRONMENT=/src/llmproxy/.venv-production
ExecStartPre=/home/noname/.local/bin/uv sync
ExecStart=/src/llmproxy/.venv-production/bin/python -m src.llmproxy.main -c /etc/llmproxy/config.yaml
Restart=on-failure
TimeoutStopSec=90s
RestartSec=5

[Install]
WantedBy=multi-user.target
