# samples/blproxy.service
[Unit]
Description=BLProxy - Backend Locking Proxy
After=network.target llama-server.service llama-reranker.service llama-embeddings.service

[Service]
Type=simple
User=noname
Group=noname
Environment=UV_PROJECT_ENVIRONMENT=/home/noname/venv-blproxy
WorkingDirectory=/src/blproxy
ExecStart=/usr/bin/env /usr/local/bin/uv run --project /src/blproxy python -m src.blproxy.main --config /src/blproxy/samples/config.yaml
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=blproxy
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

